# RaspberryPiTV/services/simpsonstv-api.service
[Unit]
Description=SimpsonsTV Control API
After=network.target simpsonstv-startup-guard.service
ConditionPathExists=!/run/tvsimpsons/skip_autostart

[Service]
Type=simple
User=root
WorkingDirectory=/home/donkikochan/TvSimpsonsApp/RaspberryPiTV
ExecStart=/usr/bin/python3 /home/donkikochan/TvSimpsonsApp/RaspberryPiTV/control_api.py
Restart=always
RestartSec=2
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
